# personal/knowledge-activation-system/scripts/com.kas.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kas.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/d/claude-code/personal/knowledge-activation-system/.venv/bin/uvicorn</string>
        <string>knowledge.api.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/d/claude-code/personal/knowledge-activation-system</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/d/claude-code/personal/knowledge-activation-system/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>KNOWLEDGE_DATABASE_URL</key>
        <string>postgresql://knowledge:localdev@localhost:5433/knowledge</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/d/logs/kas-api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/d/logs/kas-api-error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
