# personal/knowledge-activation-system/scripts/com.kas.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kas.maintenance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/d/claude-code/personal/knowledge-activation-system/.venv/bin/python</string>
        <string>/Users/d/claude-code/personal/knowledge-activation-system/web/scripts/maintenance_job.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/d/claude-code/personal/knowledge-activation-system</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/d/claude-code/personal/knowledge-activation-system/src</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/d/claude-code/personal/knowledge-activation-system/.venv/bin</string>
    </dict>

    <!-- Run daily at 3 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/kas-maintenance.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/kas-maintenance-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
